# home/Library/LaunchAgents/ke.bou.dark-mode-notify.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
"http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ke.bou.dark-mode-notify</string>
    <key>KeepAlive</key>
    <true/>
    <key>ProgramArguments</key>
    <array>
       <string>/usr/bin/swift</string>
       <string>{{ .chezmoi.homeDir }}/.config/dark-mode-notify/dark-mode-notify.swift</string>
    </array>
</dict>
</plist>
